# systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -c 10 -w 20 -f 97
